# Practice/roles/ruby_app/templates/ruby_app.service.j2
[Unit]
After=network.target

[Service]
Type=simple
WorkingDirectory={{ ruby_app_full_path }}
ExecStartPre=/bin/bash -lc '/opt/rh/{{ ruby_package }}/root/usr/local/bin/bundle exec rake db:migrate'
ExecStart=/bin/bash -lc '/opt/rh/{{ ruby_package }}/root/usr/local/bin/bundle exec puma -b unix:///var/run/puma.sock -e $RAILS_ENV config.ru'
ExecStop=/bin/kill -9 $MAINPID
Restart=on-failure
KillMode=process
{% for env in ruby_app_environment %}
Environment={{ env.name }}={{ env.value }}
{% endfor %}

[Install]
WantedBy=multi-user.target